# chip_controller_repl/rootfs/etc/services.d/ttyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start ttyd service for ingress
# ==============================================================================
bashio::log.info "Starting Web Terminal..."
cd /root || bashio::exit.nok "Can't find root folder!"

exec ttyd -p 8099 -w /root jupyter console --existing kernel.json
