# silabs-multiprotocol/rootfs/etc/services.d/socat-pty/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start socat PTY emulation
# ==============================================================================
bashio::log.info "Starting socat PTY instance..."

declare socat_options
socat_options=""

if bashio::config.true 'socat_trace'; then
    socat_options="-x -v"
fi

exec s6-notifyoncheck -s 300 -w 300 \
     "/usr/bin/socat" ${socat_options} \
         pty,raw,echo=0,link=/tmp/ttyZigbeeNCPClient,ignoreeof \
         pty,raw,echo=0,link=/tmp/ttyZigbeeNCP,ignoreeof
