# silabs-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/socat-tcp/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start socat TCP forwarder
# ==============================================================================
bashio::log.info "Starting socat TCP instance..."
exec "/usr/bin/socat" -d open:/tmp/ttyZigbeeNCPClient,nonblock,echo=0,raw,ignoreeof TCP-LISTEN:9999,reuseaddr,fork
