# api_explorer/rootfs/etc/s6-overlay/s6-rc.d/api-explorer/run
#!/command/with-contenv bashio
# shellcheck shell=bash
exec python3 /app/server.py
