# silabs-multiprotocol/rootfs/etc/services.d/zigbeed/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start EmberZNet Zigbee Daemon
# ==============================================================================
s6-svwait -U /var/run/s6/services/cpcd
s6-svwait -U /var/run/s6/services/socat-pty

mkdir -p /data/zigbeed && cd /data/zigbeed || bashio::exit.nok "Could not change to zigbeed working directory"

bashio::log.info "Starting zigbeed..."
exec "/usr/bin/stdbuf" -o0 /usr/local/bin/zigbeed
