# rt_test/rootfs/etc/profile.d/welcome.sh
#!/bin/sh

cat <<"EOF"
Use `mqtt-rt-test` to run Real-Time latency test.

EOF
