# silabs-multiprotocol/rootfs/etc/services.d/cpcd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start Co-Processor Communication Daemon (CPCd)
# ==============================================================================
bashio::log.info "Starting cpcd..."
s6-notifyoncheck -s 300 -w 300 "/usr/bin/stdbuf" -o0 /usr/local/bin/cpcd
bashio::log.info "Cleaning up after cpcd..."
rm -r /dev/shm/cpcd
