# silabs-concurrent-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/ezspd/finish
#!/usr/bin/with-contenv bashio
#==============================================================================
# EmberZNet Protocol gateway daemon finish script
#==============================================================================
bashio::log.info "ezspd ended with exit code ${1} (signal ${2})..."
