# silabs-concurrent-multiprotocol/rootfs/etc/s6-overlay/s6-rc.d/ezspd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start EmberZNet Protocol gateway deaemon
# ==============================================================================
ezsp_listen="::"
ezsp_listen_port="$(bashio::addon.port 9999)"

# If user port is not set, listen on local interface only
if ! bashio::var.has_value "${ezsp_listen_port}"; then
    ezsp_listen="$(bashio::addon.ip_address)"
    ezsp_listen_port="9999"
fi

bashio::log.info "Starting ezspd..."
exec "/usr/bin/stdbuf" -o0 \
    /usr/local/bin/ezspd
