# chip_controller/rootfs/etc/services.d/chip-repl/finish
#!/usr/bin/execlineb -S1
# ==============================================================================
# Take down the S6 supervision tree when chip-repl exits
# ==============================================================================
s6-svscanctl -t /var/run/s6/services
