# chip_controller/rootfs/etc/services.d/chip-repl/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start chip-repl service
# ==============================================================================
bashio::log.info "Starting CHIP REPL Kernel..."
cd /root || bashio::exit.nok "Can't find root folder!"

exec chip-repl
