# silabs-multiprotocol/rootfs/etc/services.d/otbr-web/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start OpenThread BorderRouter web interface
# ==============================================================================
s6-svwait -U /var/run/s6/services/otbr-agent

bashio::log.info "Starting otbr-web..."
exec /usr/sbin/otbr-web -I wpan0 -d6 -p 8080
