# chip_controller_repl/rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start ttyd service for ingress
# ==============================================================================
bashio::log.info "Starting Web Terminal..."
cd /root || bashio::exit.nok "Can't find root folder!"

declare chip_repl_options
chip_repl_options="-p /data/repl-storage.json"

if bashio::config.true 'chip_debug'; then
    chip_repl_options="${chip_repl_options} -d"
fi

exec ttyd -p 8099 -w /root chip-repl ${chip_repl_options}
