# openthread_border_router/rootfs/etc/services.d/otbr-agent/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# OpenThread BorderRouter start script
# ==============================================================================

declare otbr_agent_options
declare backbone_ifa
declare device
declare baudrate
declare flow_control
backbone_if="$(bashio::api.supervisor 'GET' '/network/info' '' '.interfaces[] | select (.primary == true) | .interface')"
otbr_agent_options="-d6 -v"
device=$(bashio::config 'device')
baudrate=$(bashio::config 'baudrate')

if bashio::config.true 'flow_control'; then
    flow_control="&uart-flow-control"
fi

if bashio::config.true 'otbr_debug'; then
    otbr_agent_options="-d7 -v"
fi

if [ -z ${backbone_if} ]; then
    bashio::log.warning "No primary network interface found! Using static eth0."
    backbone_if="eth0"
fi

mkdir -p /data/thread && ln -sft /var/lib /data/thread || bashio::exit.nok "Could not create directory /var/lib/thread to store Thread data."

bashio::log.info "Starting otbr-agent..."
exec s6-notifyoncheck -d -s 300 -w 300 -n 0 \
    "/usr/sbin/otbr-agent" -I wpan0 -B "${backbone_if}" ${otbr_agent_options} \
        "spinel+hdlc+uart://${device}?uart-baudrate=${baudrate}${flow_control}"
