# skyconnect-cp2102n-programmer/rootfs/etc/s6-overlay/s6-rc.d/skyconnect-cp2102n-programmer/up
/etc/s6-overlay/scripts/skyconnect-cp2102n-programmer.sh