# py_spy/rootfs/etc/services.d/spy-py/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start spy-py service
# ==============================================================================
declare hass_pid

# Evaluate Home Assistant PID
hass_pid=$(pgrep -f "python3 -m homeassistant --config /config")


exec py-spy record -o /share/py-spy.svg --pid ${hass_pid}