# matter_server/rootfs/etc/s6-overlay/s6-rc.d/matter-server/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start matter-server service
# ==============================================================================
bashio::log.info "Starting Matter Server..."

if bashio::config.true 'chip_debug'; then
    export CHIP_WS_DEBUG=1
fi

export CHIP_WS_STORAGE="/data/python-kv.json"
export CHIP_WS_SERVER_PORT="5580"

exec matter-server
