# chip_tool/rootfs/etc/s6-overlay/s6-rc.d/ttyd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start ttyd service for ingress
# ==============================================================================
bashio::log.info "Starting Web Terminal..."
cd /root || bashio::exit.nok "Can't find root folder!"

declare chip_repl_options
declare chip_repl_port
chip_repl_options="-p /data/repl-storage.json"
chip_repl_port="$(bashio::addon.ingress_port)"

if bashio::config.true 'chip_debug'; then
    chip_repl_options="${chip_repl_options} -d"
fi

exec ttyd --writable -p "${chip_repl_port}" -w /root tmux -u new-session -A -s homeassistant bash -l
